# systemd/nest.service
[Unit]
Description=nest daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node %h/nest/dist/cli.js start
WorkingDirectory=%h/nest
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=-%h/.config/nest/env

[Install]
WantedBy=default.target
